# inetok.service
[Unit]
Description=Internet OK

[Service]
Type=simple
ExecStart=/opt/inetok/inetok
Restart=always

[Install]
WantedBy=multi-user.target
